# rc.local
#!/bin/bash
# make ssh blocked,
/usr/sbin/sshd -D